# deploy/genilink-frontend.service
[Unit]
Description=GeniLink Frontend Service
Documentation=https://genilink.cn
After=network.target postgresql.service

[Service]
Type=simple
User=genilink
Group=genilink
WorkingDirectory=/opt/genilink-platform/frontend
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/opt/genilink-platform/frontend/.env.local

# Process management
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3001
Restart=on-failure
RestartSec=10
TimeoutStartSec=60

# Resource limits
MemoryMax=1G
CPUQuota=100%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/genilink-platform/frontend/.next-runtime
ReadWritePaths=/var/log/genilink

# Logging
StandardOutput=append:/var/log/genilink/frontend-out.log
StandardError=append:/var/log/genilink/frontend-error.log
SyslogIdentifier=genilink-frontend

[Install]
WantedBy=multi-user.target
